# deploy/yb-backup.service
[Unit]
Description=Yecny OS Nightly Backup (DB + docs)

[Service]
Type=oneshot
WorkingDirectory=/home/kruzer04/YBTM/YB-TM/yb-backend
EnvironmentFile=/home/kruzer04/YBTM/YB-TM/.env
ExecStart=/home/kruzer04/YBTM/YB-TM/venv/bin/python -m app.backup_nightly
