# deploy/yb-recurring.timer
[Unit]
Description=Run Yecny OS recurring task runner daily

[Timer]
OnCalendar=daily
Persistent=true

[Install]
WantedBy=timers.target
